# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:build]
command=npm run build